# oguri-bot.service
[Unit]
Description=OguriCap WhatsApp Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/OguriCap-Bot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10

# Variables de entorno
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PANEL_PORT=3001
Environment=SERVER_IP=178.156.179.129

# Configuración de logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oguri-bot

# Configuración de recursos
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target