# oguri-panel.service
[Unit]
Description=OguriCap Panel Next.js
After=network.target oguri-bot.service
Wants=network.target
Requires=oguri-bot.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/OguriCap-Bot/frontend-next
ExecStartPre=/usr/bin/npm install
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=15

# Variables de entorno
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=NEXT_PUBLIC_API_URL=http://178.156.179.129:3001
Environment=NEXT_PUBLIC_SOCKET_URL=http://178.156.179.129:3001

# Configuración de logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oguri-panel

# Configuración de recursos
LimitNOFILE=65536
MemoryMax=800M

[Install]
WantedBy=multi-user.target